# aws/crontabs/ins_find_similar_01.sh
#!/bin/bash
python /Users/aoran/aoranproject/manage.py ins_similar --login_account ranaoyang@outlook.com --order 1 --settings=settings.production >> /Users/aoran/aoranproject/logs/debug.log