# aws/crontabs/ins_crawl_single_thread_02.sh
#!/bin/bash
python /Users/aoran/aoranproject/manage.py ins_crawl_single_thread --api P --settings=settings.prod--uction  >> /Users/aoran/aoranproject/logs/debug.log