# aws/crontabs/ins_crawl_single_thread_01.sh
#!/bin/bash
python /Users/aoran/aoranproject/manage.py ins_crawl_single_thread --api C --settings=settings.production  >> /Users/aoran/aoranproject/logs/debug.log